# visionlink.service
[Unit]
Description=VisionLink Wearable Industrial Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=visionlink
WorkingDirectory=/home/visionlink/Desktop/visionlink
ExecStart=/usr/bin/python3 /home/visionlink/Desktop/visionlink/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
